# installer-ui/systemd/installer-ui-web.service
[Unit]
Description=DayShield Installer Web UI (busybox httpd on port 8080)
Documentation=https://github.com/daygle/dayshield-installer-ui
After=network.target
Before=installer-ui.service
ConditionKernelCommandLine=installer

[Service]
Type=simple
User=root
Group=root

# Serve /installer-ui/ as the web root.
# busybox httpd CGI: shell scripts in /installer-ui/api/ are executed directly.
# The -c flag points to the httpd config for CGI execution permissions.
WorkingDirectory=/installer-ui

# Ensure the installer has a deterministic management IP when no DHCP lease is
# available (direct-connect installs).
ExecStartPre=/bin/sh /installer-ui/bootstrap-network.sh

ExecStart=/bin/sh -c '\
  if command -v busybox >/dev/null 2>&1 && busybox --list 2>/dev/null | grep -qx httpd; then \
    if busybox httpd --help 2>&1 | grep -q -- "-t"; then \
      exec busybox httpd -f -p 0.0.0.0:8080 -h /installer-ui -c /installer-ui/httpd.conf -t 1800; \
    else \
      exec busybox httpd -f -p 0.0.0.0:8080 -h /installer-ui -c /installer-ui/httpd.conf; \
    fi; \
  elif command -v python3 >/dev/null 2>&1; then \
    exec python3 -m http.server 8080 --bind 0.0.0.0 --directory /installer-ui; \
  elif command -v python >/dev/null 2>&1; then \
    cd /installer-ui && exec python -m SimpleHTTPServer 8080; \
  else \
    echo "ERROR: No suitable HTTP server found" >&2; \
    exit 1; \
  fi'

Restart=on-failure
RestartSec=2s

# Allow LAN/browser access from installer clients.

# Security hardening - installer scripts still need to run as root
# but we restrict the file-system exposure as much as possible
PrivateNetwork=no
ProtectSystem=no
ProtectHome=yes
NoNewPrivileges=no

# Rootfs extraction (zstd | tar) can exceed 128M in CGI context.
MemoryMax=1G

[Install]
WantedBy=multi-user.target
